# xtee-misp2-base/debian/cron.d
# ID-kaartide tyhistusnimekirja uuendamine, XForms kataloogi tyhjendamine
12,42 * * * * root /usr/xtee/apache2/updatecrl.sh > /dev/null 2>&1
45 * * * * root /usr/xtee/apache2/cleanXFormsDir.sh
